# deploy/skyboard.service
# /etc/systemd/system/skyboard.service
[Unit]
Description=SkyBoard UI
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/skyboard
ExecStart=/usr/bin/npm start
Environment=NODE_ENV=production
Restart=always
User=www-data

[Install]
WantedBy=multi-user.target
